# template/deb/yaddnsc.service.in
[Unit]
Description=Yet Another Dynamic DNS Client — multi-driver DDNS updater
Requires=network-online.target
After=network.target network-online.target

[Service]
Type=simple
DynamicUser=yes
Environment=YADDNSC_CONFIG=@YADDNSC_SYSCONF_DIR@/config.json
EnvironmentFile=-@YADDNSC_SYSCONF_DIR@/default/yaddnsc
ExecStartPre=@CMAKE_INSTALL_FULL_BINDIR@/yaddnsc config -c ${YADDNSC_CONFIG} test -q
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/yaddnsc run -c ${YADDNSC_CONFIG}
Restart=on-failure
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
CapabilityBoundingSet=

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
